# notify/src/crontab
* * * * * cd /src && yt-notify >> /var/log/notify.log
* * * * * echo "the time is: $(date)" >> /var/log/cron.log
